# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
user=root
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gunicorn]
user=root
command=gunicorn wsgi:app -w 8 --user=root --group=root --chdir /app --bind=unix:/tmp/gunicorn.sock
autorestart=true
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


